# dot_config/systemd/user/homelab-backup.service.tmpl
{{- if and (eq .chezmoi.os "linux") .isServer -}}
[Unit]
Description=Backup diario del Homelab (Restic)
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/homelab/backup.env
ExecStart=%h/homelab-infra/manage.sh backup
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin

[Install]
WantedBy=default.target
{{- end -}}
